# ansible/roles/chineseroom_gitea/templates/gitea-pull-from-github.service.j2
[Unit]
Description=Pull GitHub changes into Gitea (oneshot, driven by timer)
After=gitea.service

[Service]
Type=oneshot
User={{ chineseroom_gitea_user }}
Group={{ chineseroom_gitea_group }}
ExecStart={{ chineseroom_gitea_pull_bin }} {{ chineseroom_gitea_pull_config }}
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ chineseroom_gitea_mirrors_dir }}
ReadOnlyPaths={{ chineseroom_gitea_sync_lib_dir }} {{ chineseroom_gitea_data_dir }}
