# ansible/roles/chineseroom_gitea/templates/gitea-pull-from-github.timer.j2
[Unit]
Description=Periodically pull GitHub changes into Gitea

[Timer]
OnCalendar={{ chineseroom_gitea_sync_pull_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
