# ansible/roles/chineseroom_gitea/templates/gitea-pull-from-github.service.j2
[Unit]
Description=Pull GitHub changes into Gitea (oneshot, driven by timer)
After=gitea-webhook-receiver.service

[Service]
Type=oneshot
User={{ chineseroom_gitea_user }}
Group={{ chineseroom_gitea_group }}
EnvironmentFile={{ chineseroom_gitea_webhook_env_file }}
ExecStart={{ chineseroom_gitea_pull_bin }}
NoNewPrivileges=yes
ProtectSystem=strict
