# psyopsOS/progfiguration/progfiguration/site/roles/synergycontroller/autostart/Terminal.desktop
[Desktop Entry]
Type=Application
Name=Terminal
Exec=/usr/bin/xfce4-terminal