# fly/disagreements.micahrl.com/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:remark42]
command=/srv/remark42 server
autostart=true
autorestart=true
stdout_logfile=AUTO
redirect_stderr=true
stopsignal=TERM
exitcodes=0
user=app

[program:cron]
command=crond -f
autostart=true
autorestart=true
stdout_logfile=AUTO
redirect_stderr=true
user=root
